# acmeresponder.socket
[Unit]
Description=ACME Responder
StopWhenUnneeded=true

[Socket]
ListenStream=80

[Install]
WantedBy=default.target

