# acmeresponder.service
[Unit]
Description=ACME Responder
BindsTo=dehydrated-update.service
Before=dehydrated-update.service

[Service]
ExecStart=/usr/sbin/acmeresponder
EnvironmentFile=/etc/dehydrated/config
User=dehydrated
Restart=always
StandardOutput=syslog
StandardError=syslog
