# acmeresponder.socket
[Unit]
Description=ACME Responder

[Socket]
ListenStream=80

[Install]
WantedBy=default.target

